# leaderboard_hour/leaderboard_hour_update.service
[Unit]
Description=Hourly cronjob to update leaderboard

[Service]
Type=simple
ExecStart=/opt/leaderboard_hour.py
User=ec2-user

