# leaderboard_hour/leaderboard_hour_update.timer
[Unit]
Description=script

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
Persistent=true

[Install]
WantedBy=timers.target